# ops/com.aditya.mmm.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aditya.mmm.weekly</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/aditya/Library/Mobile Documents/com~apple~CloudDocs/GitHub/MMM/scripts/run_local_workflow.sh</string>
    <string>--scheduled</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/aditya/Library/Mobile Documents/com~apple~CloudDocs/GitHub/MMM</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/aditya/Library/Mobile Documents/com~apple~CloudDocs/GitHub/MMM/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>1</integer>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/aditya/Library/Mobile Documents/com~apple~CloudDocs/GitHub/MMM/logs/launchd-weekly.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/aditya/Library/Mobile Documents/com~apple~CloudDocs/GitHub/MMM/logs/launchd-weekly.stderr.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
